# ansible/roles/image-syncthing/files/backup.service
[Unit]
Description=Backup service
After=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/syncthing-s3-backup
Restart=no
