# srvtools/ansible/roles/server/role/provision/miniserv/files/backup/backup.timer
[Unit]
Description=Backup Timer
After=network.target

[Timer]
OnCalendar=04:00
AccuracySec=10m
Persistent=true

[Install]
WantedBy=multi-user.target
